# media-sound/shairport/files/init.d/shairport
#!/sbin/runscript
# $Id$

EXEC="/usr/local/bin/shairport"
OPTIONS=""

#default
if [[ -z "$SA_USER" ]]; then SA_USER=samuel; fi
if [[ -z "$PIDFILE" ]]; then PIDFILE=/var/run/shairport.pid; fi
if [[ -z "$LOGFILE" ]]; then LOGFILE=/tmp/shairport.log; fi
if [[ -z "$ERRFILE" ]]; then ERRFILE=${LOGFILE}; fi

#configuration
if [[ ! -z "$NAME" ]]; then OPTIONS="$OPTIONS --name=$NAME"; fi
if [[ ! -z "$BUF_FILL" ]]; then OPTIONS="$OPTIONS -b $BUF_FILL"; fi
if [[ ! -z "$BACKEND" ]]; then OPTIONS="$OPTIONS --output=$BACKEND"; fi
if [[ ! -z "$PASSWORD" ]]; then OPTIONS="$OPTIONS --password=$PASSWORD"; fi
if [[ ! -z "$BACKEND_OPTS" ]]; then OPTIONS="$OPTIONS -- $BACKEND_OPTS"; fi

START_STOP_ARGS="-b -m --pidfile ${PIDFILE} -u ${SA_USER} -1 ${LOGFILE} -2 ${ERRFILE}"

start() {
	ebegin "Start shairport"
	start-stop-daemon --start \
	${START_STOP_ARGS} \
	--exec ${EXEC} \
	-- ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stop shairport"
	start-stop-daemon --stop \
	--exec ${EXEC} \
	--pidfile ${PIDFILE} 
	eend $?
}
